# pc-core-1.12.6-dev/rootfs/rootfs/etc/profile.d/pc-core.sh
#!/bin/sh

echo "PocketCluster Core version $(cat /etc/version), build $(cat /etc/pc-core)"
docker -v # "Docker version 1.3.0-dev, build ba14ddf-dirty"
